# deploy/sentiencex.service
[Unit]
Description=SentienceX-AI backend
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/sentiencex
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sentiencex/.venv/bin/gunicorn -k uvicorn.workers.UvicornWorker -w 4 -b 0.0.0.0:8000 backend.app.main:app
Restart=on-failure

[Install]
WantedBy=multi-user.target
